# sonarr.conf
[program:sonarr]
user = jobski
command = /usr/bin/sonarr -nobrowser -data=/config
